# RaspiScripts/HomeControl.target
[Unit]
Description=HomeControl Main Daemon
After=network.target
AllowIsolate=yes

[Install]
WantedBy=multi-user.target

